# synchronizer/supervisor/vtn-observer.conf
[program:computenode-observer]
command=python /opt/xos/observers/vbng/computenode-observer.py -C /opt/xos/observers/vbng/computenode_observer_config
